# TerraformWorkShop/Quadlet/SvcDisc-coredns/podman-coredns/restart.path
# ref: https://superuser.com/questions/1171751/restart-systemd-service-automatically-whenever-a-directory-changes-any-file-ins#1276457

[Unit]
Description="A unit configuration file whose name ends in \".path\" encodes information about a path monitored by systemd, for path-based activation."
Documentation=https://www.freedesktop.org/software/systemd/man/latest/systemd.path.html

[Path]
PathModified=${PathModified}

[Install]
WantedBy=default.target