# TerraformWorkShop/Quadlet/cockpit/attachments/cockpit.socket
# https://cockpit-project.org/guide/latest/listen.html#listen-systemd
[Socket]
ListenStream=127.0.0.1:9090
